# com.user.claude-launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claude-launcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/Users/jinyabo/Documents/ai-claude/claude-launcher/server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jinyabo/Documents/ai-claude/claude-launcher</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/jinyabo/Documents/ai-claude/claude-launcher/launcher.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jinyabo/Documents/ai-claude/claude-launcher/launcher.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/jinyabo</string>
  </dict>
</dict>
</plist>
